# clear-system.service
[Unit]
Description=Clear-System
[Service]
Type=oneshot
ExecStart=
[Install]
WantedBy=multi-user.target
