# clear-system.service
[Unit]
Description=Clear-System
[Service]
Type=oneshot
ExecStart=/usr/local/bin/clear-system
[Install]
WantedBy=multi-user.target
