# clear-system.timer
[Unit]
Description=My Hello world timer
 [Timer]
OnCalendar=daily
[Install]
WantedBy=timers.target
Unit=clear-system.service