# systemd/clear-system.service
[Unit]
Description=Clear-System
[Service]
Type=oneshot
ExecStart=/etc/cod-os/clear-system.sh
[Install]
WantedBy=multi-user.target
